# scripts/codeman-tunnel-named.service
[Unit]
Description=Codeman Cloudflare Named Tunnel
After=network-online.target codeman-web.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --config %h/.cloudflared/codeman.yml run codeman
Restart=always
RestartSec=5
KillMode=process

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codeman-tunnel-named

[Install]
WantedBy=default.target
